# configs/crontab/user
# Upgrade node package managers
0 3 */3 * * corepack install --global pnpm@latest > /dev/null
0 3 */3 * * corepack install --global yarn@stable > /dev/null
