# bootstrap/ubuntu/files/etc/profile.d/99-node-security.sh
# shellcheck shell=sh
export NPM_CONFIG_IGNORE_SCRIPTS=true
export npm_config_ignore_scripts=true
