# debian-ubuntu/etc/cron.daily/clear-journal-logs
#!/bin/bash

exec >/dev/null

journalctl --rotate
journalctl --vacuum-time=3d
