# debian-ubuntu/etc/cron.daily/clear-journal-logs
#!/bin/bash

journalctl --rotate
journalctl --vacuum-time=3d
