# debian-ubuntu/etc/cron.daily/clear-docker-build-cache
#!/bin/bash

exec >/dev/null

docker builder prune -f
