# debian-ubuntu/etc/cron.daily/upgrade-npm-packages
#!/bin/bash

set -e
npm upgrade -g
corepack disable pnpm
npm i bun@latest npm@latest pnpm@latest -g
npm cache clean -g --force
