# debian-ubuntu/etc/cron.daily/upgrade-npm-packages
#!/bin/bash

npm upgrade -g && corepack disable pnpm && npm i pnpm@latest -g
